# shutdownhandler.service
[Unit]
Description=Radio shutdown systemd service unit file.

[Service]
ExecStart=/shutdownhandler.out

[Install]
WantedBy=multi-user.target

